# assistant_transfer/alpine-fs/etc/profile.d/20locale.sh
version https://git-lfs.github.com/spec/v1
oid sha256:284a6ef56ab97a13a777c6b01ae14f2cc3d2b7a29c19e750e622e70cc3c73186
size 97
